# share/registration/supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
